# bot.service
[Unit]
Description=Telegram Finance Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/bot
ExecStart=/root/bot/venv/bin/python bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
